# profile.d/00_disable.el
(add-hook 'pre-command-hook (abbrev-mode nil));auto-complete.el使うからいらない
(kill-buffer "*scratch*");open-junk-fileがあるからscratchいらないです
