# profile.d/01_font.el
;;フォント設定
;;http://blog.eyesrobe.com/2007/07/22/p=35
(set-face-attribute 'default nil
		    :family "Ricty Discord"
		    :height 100)
