# profile.d/01_font.el
;;フォント設定
;;http://blog.eyesrobe.com/2007/07/22/p=35
(set-frame-font "Ricty:pixelsize=14:spacing=0")
