# profile.d/03_desktopsave.el
;;開いていたバッファをもう一度開いたり色々復元する
(desktop-save-mode 1)
