# profile.d/10_color.el
(add-to-list 'custom-theme-load-path "~/.emacs.d/bigprogram.d/emacs-color-theme-solarized/")
(load-theme 'solarized-dark t)
