# profile.d/10_close.el
(require 'undohist);;undoをファイル閉じても保存
(setq undohist-directory "~/.emacs.d/undohist")
(undohist-initialize)
