# profile.d/10_display.el
;;画像表示
(auto-image-file-mode)

;;; 現在行を目立たせる
(global-hl-line-mode)

;;; カーソルの位置が何文字目かを表示する
(column-number-mode t)

;;; カーソルの位置が何行目かを表示する
(line-number-mode t)

;;; カーソルの場所を保存する
(require 'saveplace)
(setq-default save-place t)

;;ウィンドウの上部に現在の関数名を表示します。
(which-function-mode 1)
