# profile.d/10_font.el
(global-font-lock-mode t);モードに合わせるらしい?

(set-face-attribute 'default nil
                    :family "Ricty"
                    :height 180)
