# profile.d/10_recentf.el
(custom-set-variables '(recentf-save-file "~/.emacs.d/recentf")
		      '(recentf-max-menu-items 20)
		      '(recentf-max-saved-items 500))
(require 'recentf-purge-tramp)
