# profile.d/10_recentf.el
(custom-set-variables '(recentf-save-file "~/.emacs.d/recentf"))
(require 'recentf)

(setq recentf-max-menu-items 20)
(setq recentf-max-saved-items 500)

(require 'recentf-ext)
(require 'recentf-purge-tramp)
