# profile.d/20_clipboard.el
(setq x-select-enable-clipboard t);クリップボード共有
