# profile.d/20_dired.el
;;ソート順を設定
(setq dired-listing-switches "-AFlh")
;;ごみ箱を有効 
(setq delete-by-moving-to-trash t)
