# profile.d/20_operalegendstrongest.el
;;デフォルトブラウザをOperaに
;;http://d.hatena.ne.jp/portown/20100531/1297222343
(setq browse-url-firefox-program "opera")
(setq browse-url-firefox-new-window-is-tab t)
(setq browse-url-new-window-flag t)
