# profile.d/20_recentf.el
;;(custom-set-variables '(recentf-save-file "~/Dropbox/sync/.recentf"))

(require 'recentf)
;; (require 'recentf-ext)
;; (require 'recentf-purge-tramp)
