# profile.d/20_window.el
;;http://www.bookshelf.jp/soft/meadow_30.html#SEC405
(windmove-default-keybindings);shift + arrowでウィンドウ移動
(setq windmove-wrap-around t);ループする

