# profile.d/21_key-trans.el
(define-key key-translation-map (kbd "C-t") (kbd "C-p"));悩み中
(define-key key-translation-map (kbd "M-t") (kbd "M-p"))
(define-key key-translation-map (kbd "C-M-t") (kbd "C-M-p"))
