# profile.d/30_recentf.el
;;http://d.hatena.ne.jp/rubikitch/20091224/recentf
(require 'recentf-ext)

;;http://maruta.be/intfloat_staff/101
;;メニュー表示項目
(setq recentf-max-menu-items 10000)
;;保存件数
(setq recentf-max-saved-items 100000)
