# profile.d/30_syntax-highlight.el
;;対応する括弧をハイライト
(show-paren-mode t)

;;なんかモードに合わせるらしい?よ?
(global-font-lock-mode t)

