# profile.d/50_auto-install.el
(require 'auto-install)
;;http://d.hatena.ne.jp/rubikitch/20091221/autoinstall
(setq auto-install-directory "~/.emacs.d/install/")
(auto-install-update-emacswiki-package-name t)
