# profile.d/50_disk.el
(autoload 'disk "disk" "Save, revert, or find file." t)
(global-set-key "\C-x C-s" 'disk)
