# profile.d/50_disk.el
(require 'disk)
(global-set-key (kbd "C-x C-s") 'disk)
