# profile.d/50_linum.el
;;行番号を左に表示
;;http://www23.atwiki.jp/selflearn/pages/41.html#id_51d1d876
(require 'linum)
(global-linum-mode)
