# profile.d/50_magit.el
;;emacsからgitを扱う
(add-to-list 'load-path "~/.emacs.d/bigprogram.d/magit/")
(require 'magit)
