# profile.d/50_one-key.el
;; (require 'one-key-default) ; one-key.el も一緒に読み込んでくれる
;; (require 'one-key-config) ; one-key.el をより便利にする
;; (one-key-default-setup-keys) ; one-key- で始まるメニュー使える様になる
