# profile.d/50_sudden-death.el
(require 'sudden-death)

;; ＿人人人人人人＿
;; ＞　突然の死　＜
;; ￣ＹＹＹＹＹＹ￣
