# profile.d/50_undo-tree.el
(add-to-list 'load-path "~/.emacs.d/elpa/undo-tree-0.5.5/")
(require 'undo-tree)
(global-undo-tree-mode)
