# profile.d/50_undohist.el
;;undoをファイル閉じても保存
(require 'undohist)
(undohist-initialize)
