# profile.d/50_yasnippet.el
;;http://blog.clouder.jp/archives/001037.html
;;プログラミング言語の定型文補完
(require 'yasnippet)
(yas/initialize)
