# profile.d/60_etags.el
(define-key 'helm-etags-map (kbd "C-.")		'helm-etags-select)
(define-key 'helm-etags-map (kbd "C-M-.")	'find-tag)
(define-key 'helm-etags-map (kbd "C-c .")	'visit-tags-table);;タグを再設定
(define-key 'helm-etags-map (kbd "M-.")		'pop-tag-mark)
