# profile.d/60_gdb.el
;;; 60_gdb.el --- gnu debug

;; Copyright (C) 2012  ncaq

;; Author: ncaq <ncaq@ncaq-Inspiron-N5110-LinuxMint>

(setq gdb-many-windows t)
(setq gdb-use-separate-io-buffer t) ; "IO buffer" が必要ない場合は  nil で