# profile.d/60_git.el
(require 'magit)
(setq git-commit-mode-hook nil)

(require 'git-gutter-fringe)
(global-git-gutter-mode)
