# profile.d/60_gdb.el
(setq gdb-many-windows t)
(setq gdb-use-separate-io-buffer t) ; "IO buffer"
