# profile.d/60_git.el
;; (require 'magit)

(require 'git-gutter-fringe)
(global-git-gutter-mode)
