# profile.d/60_sdic.el
(autoload 'sdic "sdic")
;;; sdic-mode 用の設定
(add-to-list 'load-path "~/.emacs.d/universe.d/sdic")
