# profile.d/60_sdic.el
(require 'sdic)
;; 辞書ファイルの設定
(setq sdic-inline-eiwa-dictionary "~/.emacs.d/universe.d/sdic/jedict.sdic")
(setq sdic-inline-waei-dictionary "~/.emacs.d/universe.d/sdic/eedict.sdic")
