# profile.d/60_sdic.el
(require 'sdic)
;;; sdic-mode 用の設定
(add-to-list 'load-path "/home/ncaq/.emacs.d/universe.d/sdic")
