# profile.d/60_slime.el
;;SLIMEのロード
(require 'slime)
(slime-setup '(slime-repl slime-fancy slime-banner))
