# profile.d/60_sdic.el
(require 'sdic)
;; 辞書ファイルの設定
(setq sdic-inline-eiwa-dictionary "/usr/local/share/dict/gene.sdic")
(setq sdic-inline-waei-dictionary "/usr/local/share/dict/jedict.sdic")
