# profile.d/70_common-lisp.el
;;CLISPをデフォルトのCommon Lisp処理系に設定
(setq inferior-lisp-program "clisp")
