# profile.d/70_haskell.el
;;Haskell関係
(require 'haskell-mode-autoloads)
(add-hook 'haskell-mode-hook 'turn-on-haskell-indentation)
