# profile.d/70_twittering-mode.el
;;http://emacswiki.org/emacs/TwitteringMode-ja
(add-to-list 'load-path "~/.emacs.d/bigprogram.d/twittering-mode/")

(require 'twittering-mode)

(setq twittering-use-master-password t)
